# universal_next_gsi/systemd/binder-proxy.service
[Unit]
Description=Android Binder Proxy Security Bridge
After=android-container.service
Requires=android-container.service
Before=mir.service

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /run/binder-proxy
# Strict allowlist proxy mapping per API 35+ capabilities.
ExecStart=/usr/bin/binder-proxy --v2 --allow camera --allow audio --allow sensors --allow vibrator --reject *
Restart=on-failure
RestartSec=3s
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
