# universal_next_gsi/systemd/mir.service
[Unit]
Description=Mir Wayland to SurfaceFlinger UI Graphics Bridge
After=binder-proxy.service
Requires=binder-proxy.service
Before=lomiri.service

[Service]
Type=notify
User=phablet
Group=video
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="MIR_SERVER_WAYLAND_HOST=wayland-0"
Environment="XDG_RUNTIME_DIR=/run/user/32011"
ExecStartPre=/usr/bin/mkdir -p /run/user/32011
ExecStartPre=/usr/bin/chown phablet:phablet /run/user/32011
# Route the startup sequence into the AVF-compatible generic bridge evaluator
# rather than hardcoding hwcomposer paths natively.
ExecStart=/system/android-bridge/gpu-bridge.sh
Restart=always
RestartSec=2s

[Install]
WantedBy=graphical.target
