# halium_gsi/systemd/hybris.service
[Unit]
Description=Initialize Libhybris Android Compatibility Layer
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target mir.service
Wants=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/mkdir -p /dev/cpuset
ExecStartPre=/bin/mount -t cgroup -o cpuset none /dev/cpuset
# Ensure linker config and Android specific props are accessible
ExecStart=/bin/sh -c "\
  export LD_LIBRARY_PATH=/system/lib64:/vendor/lib64:/system/lib:/vendor/lib;\
  # Link Android paths to the chroot\
  ln -sf /system/lib64 /lib64/android;\
  ln -sf /vendor/lib64 /lib64/vendor;\
  # Start android-binder processes if required natively\
  if [ -x /system/bin/servicemanager ]; then /system/bin/servicemanager & fi;\
"

[Install]
WantedBy=sysinit.target
WantedBy=graphical.target
