# universal_gsi/systemd/android-container.service
[Unit]
Description=Android HAL Compatibility LXC Container
Documentation=man:lxc-start(1)
After=ubuntu-rootfs.mount network.target
Before=binder-proxy.service mir.service
Wants=ubuntu-rootfs.mount

[Service]
Type=simple
# We run the container to explicitly isolate android services inside the namespace
ExecStart=/usr/bin/lxc-start -n android -f /android/container/config -F
Restart=always
RestartSec=5s

# Security Enforcement natively over container daemon
NoNewPrivileges=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
