# universal_gsi/systemd/binder-proxy.service
[Unit]
Description=Android Binder IPC Proxy (Restricted Access)
After=android-container.service
Requires=android-container.service
Before=mir.service

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /run/binder-proxy
# Assuming a proxy binary mapping only specific hal interfaces
ExecStart=/usr/bin/binder-proxy --allow camera --allow audio --allow sensors --allow vibrator
Restart=on-failure
RestartSec=3s
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
